# roles/nftables/templates/refresh-slack-ips.timer.j2
[Unit]
Description=Refresh Slack IPs every 5 minutes

[Timer]
OnBootSec=30s
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
