# roles/locksmith/templates/locksmith.service.j2
[Unit]
Description=Agent Locksmith
After=network-online.target pipelock.service
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/locksmith --config /etc/locksmith/config.yaml
EnvironmentFile=/etc/locksmith/locksmith.env
Restart=always
RestartSec=5

User=locksmith
Group=locksmith
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log/locksmith
PrivateTmp=true

[Install]
WantedBy=multi-user.target
