# roles/openclaw_service/templates/openclaw.service.j2
[Unit]
Description=OpenClaw Gateway
After=network-online.target docker.service pipelock.service
Wants=network-online.target
{% if llamafirewall.enabled | default(false) %}
After=llamafirewall.service
Wants=llamafirewall.service
{% endif %}

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/home/openclaw
ExecStart={{ final_openclaw_binary }} gateway
Restart=always
RestartSec=5

EnvironmentFile=/etc/openclaw/openclaw.env

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/home/openclaw/.openclaw
PrivateTmp=true

[Install]
WantedBy=multi-user.target
